# docker_metrics_exporter.service
[Unit]
Description=Docker Metrics Exporter for Prometheus
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
ExecStart=/usr/local/bin/docker_metrics_exporter -p 9187
Restart=on-failure
User=docker-metrics
Group=docker-metrics
# Uncomment and adjust the following line if docker is only accessible as root:
# Environment="PATH=/usr/bin:/usr/local/bin"

[Install]
WantedBy=multi-user.target


